# includes.container/etc/systemd/system/brightnesssync.service
[Unit]
Description=Zenbook Duo Brightness sync service

[Service]
ExecStart=duo watch-backlight

[Install]
WantedBy=multi-user.target